# docker_compose/samba_fs/samba_fs/root/etc/services.d/samba/run
#!/usr/bin/with-contenv bash

sleep 2

echo "smbd --foreground --no-process-group"
exec smbd --foreground --no-process-group
