# docker_compose/samba_fs/samba_fs/root/etc/services.d/winbind/run
#!/usr/bin/with-contenv bash

sleep 6

echo "winbindd --foreground --no-process-group --debuglevel=$SAMBA_FS_LOG_LEVEL"

exec winbindd --foreground --no-process-group --debuglevel=$SAMBA_FS_LOG_LEVEL
