# docker_compose/samba_dc/samba_dc/root/etc/services.d/samba_dc/run
#!/usr/bin/with-contenv bash

echo "samba -i"

nohup bash -c "/usr/local/bin/structure.sh &"

exec samba -i </dev/null