# docker_compose/samba_fs/samba_fs/root/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash

echo "avahi-daemon --no-rlimits"

exec avahi-daemon --no-rlimits
