# docker_compose/samba_dc/samba_dc/root/etc/services.d/sshd/run
#!/usr/bin/with-contenv bash

# /usr/sbin/sshd -D -ddd

nohup bash -c "/usr/local/bin/anas_zone.sh &"

exec /usr/sbin/sshd -D -E /tmp/sshd.log